# export/sync-rpi5.timer
[Unit]
Description=Sync data to rpi5 every hour

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
